# data/clr-boot-manager-booted.service.in
[Unit]
Description=clr-boot-manager tracking kernel boot success

[Service]
Type=oneshot
ExecStart=@BINDIR@/clr-boot-manager report-booted

[Install]
WantedBy=multi-user.target
