# data/clr-boot-manager-mount-boot.service.in
[Unit]
Description=mount /boot using clr-boot-manager

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@BINDIR@/clr-boot-manager mount-boot
ExecStop=@BINDIR@/umount /boot

[Install]
WantedBy=multi-user.target
